# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[program:main_backend]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app/backend
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/main_backend.log
stderr_logfile=/var/log/supervisor/main_backend.log
stdout_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1"

[program:webportal_backend]
command=python -m uvicorn main:app --host 127.0.0.1 --port 8001 --workers 1
directory=/app/webportal/backend
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/webportal.log
stderr_logfile=/var/log/supervisor/webportal.log
stdout_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1"
